# bluetooth-autoconnect.service
[Unit]
Description=Bluetooth autoconnect service
Before=bluetooth.service

[Service]
Type=simple
ExecStart=/usr/bin/bluetooth-autoconnect -d

[Install]
WantedBy=bluetooth.service
